# supervisord/conf.d/node_sonos.conf
[program:node_sonos]
command=/home/pi/.nvm/versions/node/v4.6.0/bin/npm start
directory=/home/pi/node-sonos-http-api-master
autostart=true
autorestart=true
user=pi