# supervisord/conf.d/sonos_watch.conf
[program:sonos_watch]
command=/bin/bash /home/pi/sonos_loop.sh
user=pi